# guest-scripts/xray-setup.service
[Unit]
Description=xray VM environment setup
Documentation=https://github.com/rdewit/xray
After=local-fs.target systemd-modules-load.service
Before=multi-user.target
Wants=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xray-setup.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/home

[Install]
WantedBy=multi-user.target
